# deploy/relevancy.service
[Unit]
Description=Relevancy Backend
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/relevancy/backend
Environment=PATH=/home/ubuntu/relevancy/.venv/bin
EnvironmentFile=/home/ubuntu/relevancy/backend/.env
ExecStart=/home/ubuntu/relevancy/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 2
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
